# systemd/conky.service
[Unit]
Description=Conky system monitor
After=graphical-session.target

[Service]
ExecStart=/bin/bash -c 'cd %h/.config/conky/JayBeeDe && exec /usr/bin/conky -c default.conf'
ExecStartPost=/bin/bash -c 'for i in 1 2 3 4 5; do wid=$(DISPLAY=:1 xdotool search --name "Conky" 2>/dev/null | head -1) && [ -n "$wid" ] && DISPLAY=:1 xprop -id "$wid" -remove WM_COMMAND && break || sleep 1; done'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
